# auto-trader/com.portfolio-assistant.auto-trader.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  macOS launchd plist — auto-starts the auto-trader on Mac login.

  To install:
    cp com.portfolio-assistant.auto-trader.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.portfolio-assistant.auto-trader.plist

  To uninstall:
    launchctl unload ~/Library/LaunchAgents/com.portfolio-assistant.auto-trader.plist
    rm ~/Library/LaunchAgents/com.portfolio-assistant.auto-trader.plist

  Logs: /tmp/auto-trader-stdout.log, /tmp/auto-trader-stderr.log
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.portfolio-assistant.auto-trader</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>cd ~/Git-RV/portfolio-assistant/auto-trader &amp;&amp; ./start.sh</string>
    </array>

    <!-- Start whenever user logs in — no time restriction -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Don't restart if it exits (start.sh handles its own lifecycle) -->
    <key>KeepAlive</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/tmp/auto-trader-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/auto-trader-stderr.log</string>
</dict>
</plist>
